# debian/python3-cyberfusion-queue-support.cron.d
53 8 * * * root /usr/bin/tuxis-cron /etc/cronscripts/purge-queues.cron
